# loggly/templates/default/loggly.sh.erb
#!/bin/sh -e
### BEGIN INIT INFO
# Provides:          loggly integration
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      1
# Short-Description: Register this server with loggly.
# Description:       Makes a HTTP call (via curl) to register
#                    this server as a device on loggly.
### END INIT INFO

CUSTOMER="<%=node["loggly"]["domain"]%>"
INPUT=<%=node["loggly"]["input"]%>
USER="<%=node["loggly"]["user"]%>"
PASS="<%=node["loggly"]["pass"]%>"
MASCHINE="<%=node[:scalarium][:instance][:hostname]%>"

ENDPOINT="https://${CUSTOMER}.loggly.com/api/inputs/${INPUT}"

CURL=`which curl`
CURL_OPT="-S -s -X POST -u ${USER}:${PASS}"

loggly_start () {
    local call="${CURL} ${CURL_OPT} -d name='${MASCHINE}' ${ENDPOINT}/adddevice"
    loggly_exec "$call"
}

loggly_stop () {
    local call="${CURL} ${CURL_OPT} ${ENDPOINT}/removedevice"
    loggly_exec "$call"
}

loggly_exec () {
    if [ "x$1" = "x" ]; then
        echo "Missing param."
        exit 1
    fi

    $($1)

    if [ "$?" -gt 0 ]; then
        echo "Errorz!"
        exit 1
    fi
}

case "$1" in
    start|stop)
        loggly_$1
        ;;
    *)
        echo "Usage: $0 {start|stop}"
        exit 1
        ;;
esac

exit 0
