# loggly/templates/default/loggly.sh.erb
#!/bin/sh -e
### BEGIN INIT INFO
# Provides:          loggly integration
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Register this server with loggly.
# Description:       Makes a HTTP call (via curl) to register
#                    this server as a device on loggly.
### END INIT INFO

CUSTOMER="<%=node[:loggly][:domain]%>"
INPUT=<%=node[:loggly][:input]%>
USER="<%=node[:loggly][:user]%>"
PASS="<%=node[:loggly][:pass]%>"
MASCHINE="<%=node[:scalarium][:instance][:hostname]%>"
IP="<%=node[:scalarium][:instance][:ip]%>"

ENDPOINT="https://${CUSTOMER}.loggly.com/api"

CURL=`which curl`
CURL_OPT="-o /dev/null -S -s -L -w \"%{http_code}\" -X POST -u ${USER}:${PASS}"

loggly_start () {
    local call="${CURL} ${CURL_OPT} -d ip=${IP} -d input_id=${INPUT} -d name='${MASCHINE}' ${ENDPOINT}/devices/"
    loggly_exec "$call"
}

loggly_stop () {
    local call="${CURL} ${CURL_OPT} ${ENDPOINT}/inputs/${INPUT}/removedevice/"
    loggly_exec "$call"
}

loggly_exec () {
    if [ "x$1" = "x" ]; then
        echo "Missing param."
        exit 1
    fi

    local resp=$($1)

    if [ "$?" -gt 0 ]; then
        echo "Errorz! ${resp}"
        exit 1
    fi
    echo "Registered '${MASCHINE}' with '${IP}' in loggly."
}

case "$1" in
    start|stop)
        loggly_$1
        ;;
    *)
        echo "Usage: $0 {start|stop}"
        exit 1
        ;;
esac

exit 0
