# sinopia/templates/default/sinopia.conf.erb
# Ubuntu upstart file at /etc/init/creation-service.conf
# Automatically Generated by Chef

description     "Synopia NPM cache & private repo server"

start on (starting network-interface
          or starting network-manager
          or starting networking)

stop on runlevel [!023456]

setuid <%= node['sinopia']['user'] %>
setgid <%= node['sinopia']['user'] %>

respawn

script
    exec sinopia -c <%= File.join(node['sinopia']['confdir'], 'config.yaml') %>
end script

