# easybib/tests/crontabs
0 3 * * * cd /srv/www/api/current/bin && ./getcourse cleanup changes
*/4 * * * * cd /srv/www/api/current/bin && ./getcourse resume documents
